# Configs/losgoi.service
[Unit]
Description=LineageOS Goes to IPFS

[Service]
ExecStart=/root/go/bin/LOSGoI
Restart=on-failure

[Install]
WantedBy=default.target